# roles/notifier/templates/systemd.service.j2
[Unit]
Description=Send {{ item }} notification for %i

[Service]
Type=oneshot
ExecStart=/opt/notifier/venv/bin/python /opt/notifier/notifier.py {{ item }} %i
