# roles/monitoring/files/log_srv_disk_usage.service
[Unit]
Description=Monitor srv disk usage

[Service]
Type=oneshot
User=marc
ExecStart=/bin/bash /opt/monitoring/log_srv_disk_usage.sh

