# roles/email-archive/files/mbsync-archive.timer
[Unit]
Description=Run Gmail Archive Sync Monthly

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=timers.target
